# fastcci-server.conf
description "FastCCI server upstart config"
author      "Daniel Schwen"

# used to be: start on startup
# until we found some mounts weren't ready yet while booting:
start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

script
    # We are running as user fastcci
    export HOME="/opt/fastcci"
    cd $HOME
    exec ./fastcci 8080 >> /var/log/fastcci.log 2>&1
end script

post-start script
   # Optionally put a script here that will notifiy you fastcci has (re)started
end script

